# php/crontab
# todo: only for check crontab, delete
* * * * * echo $(date) >> /proc/1/fd/2
